# deploy/systemd/mtg-leaderwebserver.service.example
[Unit]
Description=MtgLeaderwebserver
After=network.target

[Service]
Type=simple
User=mtg
Group=mtg
WorkingDirectory=/opt/mtgleaderwebserver
EnvironmentFile=/etc/mtgleaderwebserver/mtgleaderwebserver.env
# Example env file: deploy/mtgleaderwebserver.env.example
ExecStart=/opt/mtgleaderwebserver/server
Restart=on-failure
RestartSec=2

# Hardening (tune as needed)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mtgleaderwebserver

[Install]
WantedBy=multi-user.target
