# deploy/systemd/mtg-leaderwebserver.service.example
[Unit]
Description=MtgLeaderwebserver
After=network.target

[Service]
Type=simple
# Suggested install name: /etc/systemd/system/mtgleaderwebserver.service
# (the unit name you restart is the filename without the .service suffix)
User=mtg
Group=mtg
WorkingDirectory=/opt/mtgleaderwebserver
EnvironmentFile=/etc/mtgleaderwebserver/mtgleaderwebserver.env
# Example env file: deploy/mtgleaderwebserver.env.example
ExecStart=/opt/mtgleaderwebserver/server
Restart=on-failure
RestartSec=2

# Hardening (tune as needed)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mtgleaderwebserver /var/lib/mtgleader

[Install]
WantedBy=multi-user.target
